# enterprise/gateway/dingtalk-stream-bridge.service
[Unit]
Description=DingTalk Stream Bridge for OpenClaw Enterprise
After=network.target bedrock-proxy-h2.service
Wants=bedrock-proxy-h2.service

[Service]
Type=simple
User=ubuntu
Environment=DINGTALK_APP_KEY=
Environment=DINGTALK_APP_SECRET=
Environment=H2_PROXY_URL=http://127.0.0.1:8091
ExecStart=/usr/bin/python3 /opt/openclaw/enterprise/gateway/dingtalk_stream_bridge.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
